# systemd/claude-cost-tracker.service
[Unit]
Description=Aggregate Claude Code session costs into cache

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash %h/.claude/cost-tracker.sh
