# systemd/claude-cost-tracker.timer
[Unit]
Description=Run Claude cost tracker every 15 minutes

[Timer]
OnActiveSec=30s
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
